# telegram-memory-bot/telegram-memory-bot.service
[Unit]
Description=Telegram Memory Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/telegram-memory-bot
ExecStart=/opt/telegram-memory-bot/venv/bin/python bot.py
EnvironmentFile=/opt/telegram-memory-bot/.env
Restart=always
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
